# omen-thermal.service
[Unit]
Description=HP Omen 17 GPU/CPU Thermal Fan Control
After=nvidia-persistenced.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/thermal_daemon
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
